# test/fixtures/templates/stackhelper/quickstart-dotmatics/infra/infrastructor/templates/browser.service.tmpl
[Unit]
Description=browser stack
Requires=docker.service
After=docker.service

[Service]

RemainAfterExit=true
WorkingDirectory=/project/browser

ExecStartPre=

ExecStart=/usr/bin/make stack_deploy_browser ENV=production
ExecStop=/usr/bin/make stack_rm ENV=production

[Install]
WantedBy=multi-user.target
